# scripts/quizforge-ai.service
[Unit]
Description=QuizForge-AI Frontend Service
After=network.target

[Service]
Type=simple
User=leowu
WorkingDirectory=/home/leowu/Yippine/QuizForge-AI
Environment="NODE_ENV=production"
Environment="PATH=/home/leowu/.nvm/versions/node/v20.19.4/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/leowu/.nvm/versions/node/v20.19.4/bin/npm run dev
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quizforge-ai

[Install]
WantedBy=multi-user.target
